# deploy/droplet/twstockals-alert-digest-am.timer
[Unit]
Description=Daily digest AM 07:30 Taipei

[Timer]
OnCalendar=*-*-* 07:30:00
Persistent=true
Unit=twstockals-alert@digest_am.service

[Install]
WantedBy=timers.target
